# supervisor-vpl.conf
[supervisord] 
nodaemon=false

[program:vpl] 
command=/etc/init.d/vpl-jail-system start
